# www/httpd/work/.PLIST.flattened
/usr/local/man/man5/httpd.conf.5.gz
/usr/local/man/man8/httpd.8.gz
/usr/local/sbin/httpd
/usr/local/share/licenses/httpd-5.7.20150508/catalog.mk
/usr/local/share/licenses/httpd-5.7.20150508/LICENSE
/usr/local/share/licenses/httpd-5.7.20150508/BSD3CLAUSE
/usr/local/etc/rc.d/httpd
